# etc/cron.d/ksp-backup
0 * * * * ksp /bin/bash /var/lib/ksp/ksp-backup.sh | logger -t ksp
